# pbx/resources/home/django-pbx/crontab
# min hour dayofmonth month dayofweek cmd
*/10  * * * * cd /home/django-pbx/pbx; ./manage.py callstats > /dev/null 2>&1
